# templates/traefik/traefik.service.j2
[Unit]
Description=Traefik

[Service]
Type=notify
ExecStart=/usr/bin/traefik --configFile={{ traefik_config_file }}
Restart=always
WatchdogSec=1s

[Install]
WantedBy=multi-user.target
